# systemd/journal-digest.service
[Unit]
Description=Evidence Journal Daily Digest Email
After=network.target postgresql.service

[Service]
Type=oneshot
User=code
WorkingDirectory=/home/code/journal
ExecStart=/home/code/.nvm/versions/node/v24.12.0/bin/node /home/code/journal/scripts/send-digest-email.js
StandardOutput=journal
StandardError=journal
SyslogIdentifier=journal-digest

# Environment
Environment="NODE_ENV=production"
EnvironmentFile=-/home/code/journal/.env

[Install]
WantedBy=multi-user.target
