# systemd/journal-collector.service
[Unit]
Description=Evidence Journal Collector
After=network.target postgresql.service

[Service]
Type=oneshot
User=YOUR_USER
WorkingDirectory=/home/code/journal
ExecStart=/usr/bin/node src/collector/run.js
StandardOutput=journal
StandardError=journal
SyslogIdentifier=journal-collector

# Environment
Environment="NODE_ENV=production"
EnvironmentFile=-/home/code/journal/.env

[Install]
WantedBy=multi-user.target
